# utils/nonamod.conf
description "Nonamod boot script"

start on started system-services

script
    /usr/local/bin/bootmsg.sh >> /var/log/nonamod.log 2>&1
end script
